# deploy/systemd/paper-trading-daily-review.timer.example
[Unit]
Description=Run paper trading daily review at 23:55

[Timer]
OnCalendar=*-*-* 23:55:00
Persistent=true

[Install]
WantedBy=timers.target
