# deployment/runtime_dry_run/systemd/quant-alert-dry-run.service.example
# EXAMPLE ONLY — do not install directly

[Unit]
Description=Quant Alert Dry-run Runtime
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/quant-system
ExecStart=/usr/bin/python3 scripts/run_runtime_alert_center.py --dry-run
Environment=QQ_MODE=DRY_RUN
Environment=QQ_SUBMIT_PERMISSION=NO_SUBMIT
Environment=PYTHONPATH=/opt/quant-system

[Install]
WantedBy=multi-user.target
