# deployment/runtime_dry_run/systemd/quant-dashboard-render.service.example
# EXAMPLE ONLY — do not install directly

[Unit]
Description=Quant Dashboard Renderer
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/quant-system
ExecStart=/usr/bin/python3 scripts/run_runtime_operator_console.py
Environment=QQ_MODE=DRY_RUN
Environment=QQ_SUBMIT_PERMISSION=NO_SUBMIT
Environment=PYTHONPATH=/opt/quant-system

[Install]
WantedBy=multi-user.target
