# deploy/qq-bot.service
[Unit]
Description=QQ Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/qq
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/qq/.venv/bin/python3 /root/qq/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
